# services/supervisor/conf.d/ctrader_fix_live.conf
[program:ctrader_fix_live]
command=python src/GetData.py
autostart=true
autorestart=true